# src/mwan.service
[Unit]
Description=MWAN IP SLA route monitor
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/mwan/mwan.py --config /opt/mwan/mwan.toml
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
